# build/docker/root/root/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisord.sock

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && /usr/sbin/apache2 -D FOREGROUND"

[program:tika]
command=/usr/bin/java -jar /var/www/html/tika/tika-server.jar --host 127.0.0.1 --port 12345
user=www-data

[program:postgresql]
command=/usr/lib/postgresql/12/bin/postgres -D /var/www/postgresql
user=www-data

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server

